# supervisor/conf.d/ssserver.conf
[program:ssserver]
command=/root/fig/ssserver/ssserver.sh
